# crontab
*/1 * * * * /usr/bin/python2.7  /mnt/sda1/arduino/gardencode/plotSoilSensorSingleShot.py
